# systemd/clawbie-chat-bridge.timer
[Unit]
Description=Run Clawbie OpenClaw chat-memory bridge every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
RandomizedDelaySec=30s
Unit=clawbie-chat-bridge.service

[Install]
WantedBy=timers.target
