# systemd/clawbie-sub-agent-activity.timer
[Unit]
Description=Track Codex sub-agent activity every 10 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
Persistent=true
Unit=clawbie-sub-agent-activity.service

[Install]
WantedBy=timers.target

