# systemd/clawbie-nightly-reverie.service
[Unit]
Description=Nightly Reverie Generation
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/mal/git/clawbie
ExecStart=/home/mal/git/clawbie/skills/nightly-reverie/scripts/run_direct_reverie.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target